# deploy/systemd/buzz-news-worker.service
[Unit]
Description=BuzzNews background worker
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/buzznews
EnvironmentFile=/home/ubuntu/buzznews/.env
ExecStart=/home/ubuntu/buzznews/.venv/bin/python -m buzz_news run-worker
Restart=always
RestartSec=5
StandardOutput=append:/var/log/buzz-news/worker.log
StandardError=append:/var/log/buzz-news/worker.err.log

MemoryMax=350M
CPUQuota=150%

[Install]
WantedBy=multi-user.target
